# config/baseball-dashboard.service
# Legacy/copy service unit kept for reference during cleanup.
# Canonical repo service unit source is: web/college-baseball-dashboard.service
# Do not assume this file matches the active deployed unit without diffing first.
[Unit]
Description=College Baseball Predictor Dashboard
After=network.target

[Service]
Type=simple
User=sam
WorkingDirectory=/home/sam/college-baseball-predictor
ExecStart=/usr/bin/python3 /home/sam/college-baseball-predictor/web/app.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
